# cluster-dist-shell/reinstall/rc.local
#!/bin/bash

mount sfs1:/sharedfs01/home /home
mount sfs1:/sharedfs01/opt.shared /opt/shared
